# .config/systemd/user/cls.timer
[Unit]
Description=Update checker for cls

[Timer]
OnUnitActiveSec=10min
OnBootSec=10min

[Install]
WantedBy=timers.target
